# rootfs/etc/s6-overlay/s6-rc.d/mega-serve/run
#!/usr/bin/env sh

setsid /usr/local/bin/mega-cmd-server
