# rootfs/etc/s6-overlay/s6-rc.d/mega-login/up
#!/usr/bin/env sh

: ${EMAIL}
: ${PASSWORD}

sleep 1

/usr/local/bin/mega-exec login "$EMAIL" "$PASSWORD"
