# rootfs/etc/s6-overlay/s6-rc.d/mega-login/up
#!/usr/bin/env sh

: ${MEGA_EMAIL}
: ${MEGA_PASSWORD}

sleep 1

/usr/local/bin/mega-exec login "$MEGA_EMAIL" "$MEGA_PASSWORD"
